# hys_relay.conf
description "Set permissions for halfyserial relay access"

start on (filesystem)

task

exec /var/www/halfyserial/io_init
